# assets/s6-overlay/services/postgres-ready/up
#!/usr/bin/execlineb 

# foreground {
#   redirfd -w 1 /dev/null
#   fdmove -c 2 1
#   s6-setuidgid postgres
#   exec /usr/lib/postgresql/13/bin/pg_ctl -D /datadir/postgresql/data stop 2>&1
# }

/etc/s6-overlay/script/postgres-ready