# utils/s6-overlay/services/postgres-ready/up
#!/usr/bin/execlineb 

/etc/s6-overlay/script/postgres-ready